# n0stack-backup.service
[Unit]
Description=n0stack backup daemon. specify ip address.

[Service]
ExecStart=/fio/n0stack-backup/backup.sh
Type=simple
