# n0stack-backup-image.timer
[Unit]
Description=n0stack image backup

[Timer]
OnCalendar=*:00
Persistent=true
Unit=n0stack-backup-image.service

[Install]
WantedBy=timer.target
