# n0stack-backup-image.service
[Unit]
Description=n0stack image backup

[Service]
ExecStart=/fio/n0stack-backup/backup.sh
Type=simple
